# roles/munin/files/munin-cgi-graph.socket
[Unit]
Description=Munin socket

[Socket]
ListenStream=/var/run/munin/fcgi-munin-graph.sock

[Install]
WantedBy=sockets.target
