# roles/rtorrent/files/rtorrent.service.j2
[Unit]
Description=rTorrent
After=network.target

[Service]
ExecStart=rtorrent
User={{ username }}

[Install]
WantedBy=multi-user.target
